# ansible/roles/postgres/templates/pgbackrest-diff.timer.j2
[Unit]
Description=Nightly pgBackRest differential backup (Mon-Sat)

[Timer]
OnCalendar=Mon..Sat *-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
