# ansible/roles/postgres/templates/pgbackrest-full.timer.j2
[Unit]
Description=Weekly pgBackRest full backup (Sunday)

[Timer]
OnCalendar=Sun *-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
